# Node/Aztec/aztec-node.service
[Unit]
Description=Aztec Node Service
After=network.target docker.service
Requires=docker.service

[Service]
User=root
WorkingDirectory=/home/root
ExecStart=/home/root/start_aztec_node.sh
Restart=always
RestartSec=10
Environment=PATH=/usr/bin:/bin:/usr/local/bin:/home/root/.aztec/bin

[Install]
WantedBy=multi-user.target
